# device/dns_monitor.service
[Unit]
Description=AI Privacy Firewall DNS Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ai-privacy-firewall/device
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/pi/ai-privacy-firewall/device/dns_monitor.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
